# Superiori/Documentazioni/TomCat/Files/rc.local
#! /bin/bash
/opt/tomcat/bin/startup.sh
/home/enrypase/Desktop/startup.sh
exit 0
